# data/hp-helperd.service
[Unit]
Description=HP Helper privileged hardware read daemon (Python)
After=systemd-udevd.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m hp_helperd
WorkingDirectory=@ROOT_DIR@
Environment=PYTHONPATH=@ROOT_DIR@
Restart=on-failure
RestartSec=1
RuntimeDirectory=hp-helperd
RuntimeDirectoryMode=0755

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/run/hp-helperd -/sys/class/hwmon -/sys/class/leds -/sys/devices/platform/hp-wmi -/sys/devices/platform/hp-kbd-rgb
ReadOnlyPaths=/sys/class/powercap /usr/local/bin/ryzenadj
RestrictAddressFamilies=AF_UNIX
DeviceAllow=/dev/mem r
DeviceAllow=char-input r

[Install]
WantedBy=multi-user.target
